# server/ap-fallback.service
[Unit]
Description=Switch Charging V2 - AP fallback when no client WiFi at boot
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/pi/Code/Switch_Charging_V2
ExecStart=/usr/bin/python3 -m server.ap_fallback
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
